# templates/default/mcrouter.conf.erb
# This file is managed by Chef, Any Changes will be overwritten
# mcrouter - memcached protocol router
#
# mcrouter is a memcached protocol router for
# scaling memcached deployments

description "memcached protocol router"

start on runlevel [2345]
stop on runlevel [2345]

respawn

setuid <%= node['mcrouter']['user'] %>

exec /usr/local/bin/mcrouter <%= @cli_opts %>
